# templates/etc/supervisor/conf.d/nsqadmin.conf.j2
[program:nsqadmin]

command={{ nsq_nsqadmin_cmd }} -config=/etc/nsqadmin.conf

stdout_logfile=/var/log/nsqadmin.log
stderr_logfile=/var/log/nsqadmin.log
