# templates/etc/supervisor/conf.d/nsqd.conf.j2
[program:nsqd]

user={{ nsq_nsqd_user }}

command={{ nsq_nsqd_cmd }} -config=/etc/nsqd.conf

stdout_logfile=/var/log/nsqd.log
stderr_logfile=/var/log/nsqd.log
