# templates/etc/supervisor/conf.d/nsqlookupd.conf.j2
[program:nsqlookupd]

command={{ nsq_nsqlookupd_cmd }} -config=/etc/nsqlookupd.conf

stdout_logfile=/var/log/nsqlookupd.log
stderr_logfile=/var/log/nsqlookupd.log
